# misc/devenv/crontab
* * * * * /usr/libexec/mcollective/update_yaml.rb
