# misc/devenv/crontab
* * * * * /usr/libexec/mcollective/update_yaml.rb /etc/mcollective/facts.yaml &> /dev/null
0 * * * * /usr/sbin/oo-last-access > /var/lib/openshift/last_access.log 2>&1
