# docs/devenv/crontab
* * * * * /usr/bin/puppet /usr/libexec/mcollective/update_yaml.pp
